# quadlets/pwpush-ephemeral.container
[Container]
ContainerName=pwpush
Image=docker.io/pglombardo/pwpush:latest
PublishPort={{HOST_PORT}}:5100

[Service]
Restart=always

[Install]
WantedBy=default.target